# infrastructure/ansible/roles/rpi_prep/templates/rpi-throttling-textfile-exporter.service.j2
[Unit]
Description=Export Raspberry Pi throttling state for node-exporter
Documentation=man:vcgencmd(1)

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/rpi-throttling-textfile-exporter
